# battery-script.service
[Unit]


[Service]
ExecStart=/usr/local/bin/battery-script.sh

[Install]
WantedBy=default.target

