# battery-script.service
[Unit]


[Service]
ExecStart=/usr/local/bin/battery.sh
ExecStop=/usr/local/bin/battery-stop.sh

[Install]
WantedBy=default.target

